# systemd/1ai-reach-dashboard.service
[Unit]
Description=1ai-reach Next.js Dashboard
After=network.target 1ai-reach-mcp.service

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/projects/1ai-reach/dashboard
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/Cellar/node@24/24.14.1/bin:/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
ExecStart=/home/linuxbrew/.linuxbrew/Cellar/node@24/24.14.1/bin/node /home/openclaw/projects/1ai-reach/dashboard/node_modules/.bin/next start -p 8502
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
