# systemd/1ai-reach-ui.service
[Unit]
Description=1ai-reach Streamlit UI
After=network.target 1ai-reach-mcp.service

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/.openclaw/workspace/1ai-reach
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/linuxbrew/.linuxbrew/bin/streamlit run ui/app.py --server.port 8502 --server.address 0.0.0.0 --server.headless true
Restart=always
RestartSec=10
StandardOutput=append:/home/openclaw/.openclaw/workspace/1ai-reach/logs/streamlit.log
StandardError=append:/home/openclaw/.openclaw/workspace/1ai-reach/logs/streamlit.log

[Install]
WantedBy=multi-user.target
